# contour-update.timer
[Unit]
Description=Run contour-app auto-update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

